# realtime/scripts/trading-engine.service
[Unit]
Description=Real-time Trading Engine
After=network.target postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/mnt/c/Users/emily/Documents/GitHub/get-money-get-paid/realtime
Environment=ALPACA_API_KEY=your_api_key_here
Environment=ALPACA_SECRET_KEY=your_secret_key_here
ExecStart=/opt/dagster/venv/bin/python -m src.main --log-level INFO
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
